# pansnaps.crontab
MAILTO=""
* * * * * $OMD_ROOT/etc/thruk/plugins-enabled/bin/pansnaps
